# systemd/kokkan.timer
[Unit]
Description=Run Kokkan every 10 seconds

[Timer]
OnBootSec=10
OnUnitActiveSec=10
Unit=kokkan.service

[Install]
WantedBy=timers.target
